# root_fs/os/etc/crontab
*/5 * * * * /usr/sbin/logrotate /var/logrotate.d/sys_logrotate.conf 
