# root_fs/os/usr/lib/opkg/info/br2684ctl.list
/usr/sbin/br2684ctl
/etc/init.d/br2684ctl
