# root_fs/os/usr/lib/opkg/info/lxc.list
/usr/bin/lxc-unshare
/usr/share/lxc/selinux/lxc.if
/usr/share/lxc/config/ubuntu-cloud.lucid.conf
/usr/lib/lxc/lxc-net
/usr/bin/lxc-freeze
/usr/bin/lxc-unfreeze
/usr/share/lxc/config/opensuse.userns.conf
/usr/lib/lxc/lxc-user-nic
/usr/share/lxc/templates/lxc-oracle
/usr/share/lxc/config/debian.userns.conf
/usr/bin/lxc-ls
/usr/share/lxc/config/ubuntu.userns.conf
/usr/share/lxc/templates/lxc-plamo
/usr/bin/lxc-checkconfig
/usr/bin/lxc-create
/usr/share/lxc/config/centos.common.conf
/usr/share/lxc/templates/lxc-busybox
/usr/share/lxc/hooks/squid-deb-proxy-client
/lib/systemd/system/lxc-net.service
/lib/liblxc.so.1.2.0
/usr/share/lxc/config/gentoo.userns.conf
/usr/bin/lxc-device
/etc/init/lxc-net.conf
/usr/bin/lxc-copy
/lib/liblxc.so
/lib/pkgconfig/lxc.pc
/usr/share/lxc/hooks/clonehostname
/usr/share/lxc/templates/lxc-fedora
/usr/bin/lxc-monitor
/usr/bin/lxc-start
/lib/liblxc.so.1
/usr/share/lxc/config/ubuntu-cloud.common.conf
/usr/share/lxc/config/fedora.common.conf
/usr/share/lxc/config/slackware.common.conf
/usr/share/lxc/config/sparclinux.userns.conf
/usr/share/lxc/config/openwrt.common.conf
/usr/share/lxc/config/common.conf
/usr/bin/lxc-top
/usr/share/lxc/hooks/mountecryptfsroot
/etc/init/lxc.conf
/etc/init.d/lxc.sh
/usr/share/lxc/config/archlinux.common.conf
/usr/share/lxc/selinux/lxc.te
/usr/bin/lxc-info
/usr/share/lxc/config/nesting.conf
/etc/lxc/default.conf
/usr/share/lxc/templates/lxc-sshd
/usr/bin/lxc-checkpoint
/lib/systemd/system/lxc@.service
/usr/share/lxc/templates/lxc-ubuntu-cloud
/usr/share/lxc/lxc-patch.py
/usr/share/lxc/templates/lxc-sparclinux
/usr/share/lxc/config/oracle.userns.conf
/include/lxc/attach_options.h
/etc/init/lxc-instance.conf
/usr/lib/lxc/lxc-containers
/etc/rc3.d/S64lxc
/usr/share/lxc/config/alpine.userns.conf
/etc/default/lxc
/usr/share/lxc/templates/lxc-openmandriva
/usr/share/lxc/config/gentoo.moresecure.conf
/usr/share/lxc/config/ubuntu.lucid.conf
/include/lxc/lxccontainer.h
/usr/share/lxc/config/plamo.common.conf
/usr/bin/lxc-destroy
/usr/bin/lxc-config
/lib/lua/5.1/lxc/core.so
/usr/share/lxc/config/centos.userns.conf
/usr/share/lxc/templates/lxc-gentoo
/usr/share/lxc/config/opensuse.common.conf
/usr/share/lxc/templates/lxc-download
/usr/share/lxc/config/debian.common.conf
/usr/share/lua/5.1/lxc.lua
/usr/bin/lxc-usernsexec
/usr/share/lxc/config/ubuntu.common.conf
/usr/bin/lxc-execute
/usr/share/lxc/templates/lxc-archlinux
/usr/share/lxc/templates/lxc-cirros
/usr/sbin/init.lxc
/usr/share/lxc/config/ubuntu-cloud.userns.conf
/usr/bin/lxc-snapshot
/usr/share/lxc/config/fedora.userns.conf
/usr/share/lxc/config/gentoo.common.conf
/usr/share/lxc/config/slackware.userns.conf
/usr/bin/lxc-autostart
/usr/share/lxc/lxc.functions
/usr/share/lxc/config/common.conf.d/README
/usr/share/lxc/config/userns.conf
/usr/share/lxc/config/archlinux.userns.conf
/usr/share/lxc/templates/lxc-debian
/usr/share/lxc/config/sparclinux.common.conf
/usr/share/lxc/hooks/ubuntu-cloud-prep
/usr/lib/lxc/hooks/unmount-namespace
/usr/share/lxc/templates/lxc-ubuntu
/usr/lib/lxc/lxc-apparmor-load
/usr/bin/lxc-cgroup
/lib/lxc/rootfs/README
/usr/bin/lxc-console
/usr/lib/lxc/lxc-monitord
/usr/share/lxc/config/common.seccomp
/usr/bin/lxc-attach
/lib/systemd/system/lxc.service
/usr/share/lxc/templates/lxc-centos
/usr/share/lxc/templates/lxc-altlinux
/include/lxc/version.h
/usr/share/lxc/config/oracle.common.conf
/usr/share/lxc/templates/lxc-alpine
/usr/share/lxc/templates/lxc-slackware
/usr/bin/lxc-wait
/usr/share/lxc/config/plamo.userns.conf
/usr/share/lxc/templates/lxc-opensuse
/usr/share/lxc/config/alpine.common.conf
/usr/bin/lxc-stop
