# root_fs/os/usr/lib/opkg/info/zywifid.list
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_5G_CE_ODM.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_SingleSKU_BF_20181002.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20B/SingleSKU_Shatel_CE.dat
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_2.4G_CE_ODM.dat
/usr/ecn/ECN_PowerTable/MT7915A/SingleSKU_Table/001_DX3401-T0_TMobile.dat
/usr/ecn/ECN_PowerTable/EMG3525-T50B/7615_SingleSKU_BF_CE_20190619.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_BF_CE_20190313.dat
/usr/ecn/ECN_PowerTable/PMG5817-T20A/7615_SingleSKU_BF_20190820.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_SingleSKU_20180821.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_CE_20190418.dat
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_2.4G_BR_Anatel_ODM.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/009_HGW500TX2X2E-T0_CE_Mandy.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/009_HGW500TX2X2E-T0_CE_BF_Mandy.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20A/SingleSKU_MAX.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/008_DX3301-T0_TelenorDigi_MY_BF_Summer.dat
/etc/rc3.d/S55zywifid
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/008_DX3301-T0_TelenorDigi_MY_Summer.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20B/SingleSKU_Turkey.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20B/SingleSKU_Internal_antenna_EU_20190510.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_Single_SKU_BF_20181113.dat
/usr/ecn/ECN_PowerTable/VMG3625-T50B/7615_SingleSKU_BO_20200429.dat
/bin/bcm_erp.sh
/usr/ecn/ECN_PowerTable/VMG3927-T50K/7615_SingleSKU_BF_PH_20210422.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/004_EX3301-T0_TelecomMalaysia_CE_Summer.dat
/usr/ecn/ECN_PowerTable/VMG1302-T20A/SingleSKU_CE.dat
/usr/ecn/ECN_PowerTable/PMG5317-T20A/SingleSKU_TW_5dB.dat
/usr/ecn/ECN_PowerTable/EMG3525-T50B/7615_SingleSKU_CE_20190619.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/007_WX3100-T0_TelenorDigi_MY_BF_Emma.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20B/SingleSKU_Shatel_MAX.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/011_EX2210-T0_Backoff_Table_FCC_20220426.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_CE_20190313.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_SingleSKU_TIM_20180711.dat
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_5G_BR_Anatel_ODM.dat
/usr/ecn/ECN_PowerTable/EMG5523-T50B/7615_SingleSKU_BF_MY_20201113.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/003_WX3100-T0_CE_BF_Emma.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/004_EX3301-T0_TelecomMalaysia_CE_BF_Summer.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20A/SingleSKU_CE_0930.dat
/usr/ecn/ECN_PowerTable/VMG1302-T20A/SingleSKU_MAX.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_Single_SKU_BF_Table_20180911.dat
/bin/wlan_wps
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_External_antenna_5G_CE_ODM.dat
/usr/ecn/ECN_PowerTable/VMG3927-T50K/7615_SingleSKU_PH_20210422.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20B/SingleSKU_Internal_antenna_20190212.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/002_DX3301-T0_FCC_BF.dat
/usr/ecn/ECN_PowerTable/EMG3525-T50B/7615_Single_SKU_BF_CHT_20191203.dat
/usr/ecn/ECN_PowerTable/EMG3525-T50B/7615_Single_SKU_CHT_20191203.dat
/usr/lib/wifi_module/libzywl_001_ecn.so
/bin/wifi_config_save
/usr/ecn/ECN_PowerTable/EMG5523-T50B/7615_SingleSKU_MY_20201113.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/001_DX3301-T0_Generic_CE_Summer.dat
/usr/ecn/ECN_PowerTable/VMG3312-T20A/SingleSKU_CE.dat
/bin/map_start
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/006_EX3200-T0_Generic_CE_BF_Emma.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_SingleSKU_TIM_20180516.dat
/bin/zywifid
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_BF_PH_20210701.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_SingleSKU_20181002.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/003_WX3100-T0_CE_Emma.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_BF_RU_20220214.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50A/7615_SingleSKU_20181218.dat
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_5G_CE_Lab.dat
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_2.4G_CE_Lab.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50A/7615_SingleSKU_BF_20181218.dat
/bin/zywifid_run.sh
/usr/ecn/ECN_PowerTable/VMG3625-T50B/7615_SingleSKU_BF_BO_20200429.dat
/usr/ecn/ECN_Region_MT7915D.txt
/usr/ecn/ECN_PowerTable/PMG5817-T20A/7615_SingleSKU_20190820.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20A/SingleSKU_Anatel.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20A/SingleSKU_CE_Final_1212.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/006_EX3200-T0_Generic_CE_Emma.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/007_WX3100-T0_TelenorDigi_MY_Emma.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_BF_CA_20190619.dat
/usr/ecn/ECN_PowerTable/VMG3312-T20A/SingleSKU_MAX.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/001_DX3301-T0_Generic_CE_BF_Summer.dat
/usr/ecn/ECN_PowerTable/VMG3625-T20A/SingleSKU_External_antenna_2.4G_CE_ODM.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_PH_20210701.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_SingleSKU_TIM_20180724.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_RU_20220214.dat
/bin/wifi_config_save_receive
/usr/ecn/ECN_PowerTable/PMG5317-T20A/SingleSKU_TW_3dB.dat
/etc/init.d/zywifid.sh
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/002_DX3301-T0_FCC.dat
/usr/ecn/ECN_PowerTable/VMG8825-T50A/7615_Single_SKU_20181113.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/010_EX3300-T0_MT7915_Backoff_Table_CHT_20210824_Summer.dat
/usr/ecn/ECN_PowerTable/PMG5317-T20B/SingleSKU_2.4G_CE_ODM.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_CE_20190827.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_BF_CE_20190418.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/011_EX2210-T0_MT7915_Sku_Table_FCC_20220406.dat
/usr/ecn/ECN_PowerTable/MT7915A/SingleSKU_BF_Table/001_DX3401-T0_TMobile_BF.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_BF_Table/005_WX3100-T0_CE_BF_Emma.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_FCC_20200723.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_BF_FCC_20200723.dat
/usr/ecn/ECN_PowerTable/VMG1312-T20A/SingleSKU_CE.dat
/usr/ecn/ECN_PowerTable/VMG8623-T50B/7615_SingleSKU_CA_20190619.dat
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/010_EX3300-T0_MT7915_Sku_Table_CHT20210824_Summer.dat
/bin/zywlctl
/usr/ecn/ECN_PowerTable/MT7915D/SingleSKU_Table/005_WX3100-T0_CE_Emma.dat
