# rootfs/os/usr/lib/opkg/info/base-files-econet.list
/lib/modules/slic3_sep.ko
/lib/modules/ksocket.ko
/etc/init.d/makeDevs.sh
/etc/resolv.conf
/tmp/.gitkeep
/proc/.gitkeep
/etc/preinit
/etc/checkXdsl.sh
/misc/.gitkeep
/mnt/.gitkeep
/etc/init.d/button_disable_for_manufacture_extender.sh
/usr/etc/map_cfg_agent.txt
/etc/rtspStream.sh
/lib/modules/fxs3_silicon.ko
/usr/etc/wts_bss_info_config
/usr/etc/mapd_cfg_agent.txt
/lib/modules/lec.ko
/lib/modules/pcm1.ko
/lib/modules/ortp.ko
/lib/modules/spi.ko
/lib/upgrade/keep.d/base-files-econet
/etc/quarterHour.sh
/etc/fstab
/usr/etc/passwd
/usr/etc/init.d/rcS
/etc/init.d/zyusbinit.sh
/sys/.gitkeep
/lib/modules/pcm2.ko
/lib/modules/pcmDump.ko
/etc/initQuarterHourFiles.sh
/lib/preinit/makeDevNodes.sh
/usr/etc/protocols
/etc/isolateDMZ_VMG8825_T50K.sh
/run
/userfs/bin/vconfig
/usr/etc/inittab
/usr/etc/map_cfg.txt
/usr/etc/profile
/usr/etc/services
/usr/etc/shells
/lib/modules/acodec_x.ko
/var/.gitkeep
/lib/modules/fxs3.ko
/usr/etc/mapd_strng.conf
/lib/preinit/preinit
/lib/modules/slic3.ko
/etc/init.d/button_disable_for_manufacture.sh
/data/.gitkeep
/lib/modules/ovdsp.ko
/lib/modules/foip.ko
/dev/.gitkeep
/usr/etc/banner
/userfs/scripts/start_lan.sh
/usr/script/mesh_start.sh
/usr/etc/mtab
/lib/modules/slic3_silicon.ko
/etc/init.d/rcS
/usr/etc/mapd_cfg.txt
/etc/inittab
/lib/preinit/makeExtDevNodes.sh
/usr/etc/group
/etc/init.d/miniboot_rcS
/usr/etc/hosts
/etc/miniboot.sh
/lib/modules/sys_mod.ko
/lib/modules/DSPCore.ko
