# rootfs/os/usr/lib/opkg/info/busybox.list
/bin/dmesg
/usr/bin/id
/bin/login
/bin/rm
/usr/bin/wget
/sbin/lsmod
/bin/sh
/usr/bin/lsof
/sbin/vconfig
/usr/bin/test
/sbin/depmod
/bin/zcat
/bin/fgrep
/sbin/halt
/bin/ash
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/sbin/logread
/usr/sbin/ntpd
/bin/stty
/usr/bin/passwd
/sbin/klogd
/usr/bin/traceroute
/bin/more
/usr/bin/sendarp
/bin/hostname
/bin/vi
/usr/bin/tftp
/sbin/modprobe
/bin/pwd
/usr/bin/nc
/bin/chown
/etc/init.d/cron
/bin/mount
/bin/sync
/usr/bin/dirname
/sbin/rmmod
/usr/bin/hostid
/sbin/hwclock
/bin/bash
/sbin/insmod
/usr/bin/traceroute6
/bin/ping6
/usr/bin/basename
/bin/dnsdomainname
/bin/echo
/usr/bin/less
/usr/sbin/chroot
/usr/sbin/crond
/usr/bin/head
/bin/pingsvrs
/bin/gunzip
/bin/pings
/bin/netstat
/bin/busybox
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/sbin/watchdog
/bin/kill
/bin/mktemp
/usr/bin/tr
/bin/mknod
/bin/rmdir
/usr/bin/tee
/usr/bin/ftpput
/usr/bin/bzcat
/usr/bin/hexdump
/bin/date
/usr/bin/wc
/usr/bin/nslookup
/usr/bin/[
/usr/bin/pgrep
/bin/umount
/usr/bin/tail
/usr/bin/eject
/bin/mpstat
/usr/bin/mkfifo
/bin/chgrp
/bin/cp
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/df
/bin/false
/usr/bin/killall
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/usr/bin/cmp
/usr/bin/reset
/sbin/start-stop-daemon
/bin/usleep
/bin/cat
/etc/init.d/telnet
/usr/bin/arping
/usr/bin/telnet
/sbin/getty
/bin/touch
/bin/tar
/usr/bin/taskset
/usr/bin/yes
/usr/bin/top
/sbin/switch_root
/usr/bin/expr
/usr/bin/env
/sbin/sysctl
/sbin/syslogd
/bin/sed
/usr/bin/logger
/bin/chmod
/usr/bin/awk
/bin/uname
/usr/bin/which
/usr/bin/find
/bin/ln
/bin/ls
/sbin/pivot_root
/usr/bin/cut
/usr/bin/du
/bin/mv
/usr/sbin/telnetd
/usr/bin/crontab
/usr/bin/printf
/sbin/init
/usr/bin/tty
/bin/ping
/usr/bin/bunzip2
/usr/sbin/chpasswd
/usr/bin/uptime
/bin/fsync
/bin/pidof
/bin/egrep
/usr/bin/readlink
/sbin/arp
/bin/ps
/usr/bin/xargs
/bin/gzip
