# rootfs/os/usr/lib/opkg/info/igmpproxy.list
/etc/init.d/igmpproxy
/bin/igmpproxy
/etc/config/igmpproxy
/etc/igmpproxy.conf
