# rootfs/os/usr/lib/opkg/info/zcfg_be.list
/etc/init.d/smbd.sh
/etc/init.d/radvd.sh
/lib/private/libzcfg_be.so
