# rootfs/os/usr/lib/opkg/info/zyshcmd_wrapper.list
/etc/init.d/zyshcmd_wrapper.sh
/etc/rc3.d/S51zyshcmd_wrapper
/bin/zyshcmd_wrapper
