# root_fs/etc/init.d/save_matrix
#!/bin/sh /etc/rc.common
START=99
STOP=01

start() {
    exit 0
}

stop() {
    /sbin/save_matrix
}
